# SendMail.conf
start on runlevel [2345]
stop on runlevl [!2345]
respawn
exec /home/SendMail/SendMail.sh
